# supervisord.conf
[supervisord]
nodaemon=true

[program:ollama]
command=ollama serve
autostart=true
autorestart=true

[program:fastapi]
command=python main.py
autostart=true
autorestart=true
